# sensor-node/smartsense-sensor.service
[Unit]
Description=SmartSense Sensor Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/smartsense/sensor-node
ExecStart=/usr/bin/python3 /home/pi/smartsense/sensor-node/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
